# scripts/roscar-cv.service
[Unit]
Description=ROScar1 CV Stack (ArUco + YOLO on WSL2 with GPU)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=brian
Environment=HOME=/home/brian
WorkingDirectory=/home/brian
ExecStart=/home/brian/ROScar1/scripts/roscar-cv.sh
Restart=always
RestartSec=5
# Don't restart-loop forever if something is fundamentally broken
StartLimitIntervalSec=120
StartLimitBurst=5

# Give ROS2 nodes time to shut down gracefully
TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGINT

# Logging to journalctl
StandardOutput=journal
StandardError=journal
SyslogIdentifier=roscar-cv

[Install]
WantedBy=multi-user.target
